# config/systemd/user/agent-maintenance.timer
[Unit]
Description=Run daily agent maintenance

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
