# config/systemd/user/agent-worktree-prune.timer
[Unit]
Description=Run weekly agent worktree prune

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
