# config/systemd/user/tmux-snapshot.timer
[Unit]
Description=Periodic tmux agent registry snapshot

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
